# systemd/subway-sign.service.example
[Unit]
Description=NYC MTA Subway Sign Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/YOUR_USERNAME/subway-sign
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/subway-sign/run.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Performance tuning
Nice=-10
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
